# templates/etc/systemd/system/docker.service.d/options.conf.j2
# {{ ansible_managed }}

[Service]
ExecStart=/usr/bin/dockerd {{ docker_daemon_flags | join(" ") }}